# ops/projectpulse/systemd/projectpulse-restore-validation.service
[Unit]
Description=Project Health Dashboard Backup Restore Validation Export

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/projectpulse-restore-validation-export.sh
